# crontab
* * * * * /usr/local/bin/php /var/www/html/cron.php >> /var/www/html/cron_output.txt 2>&1



